# data/battery-input-manager_user.service.in
[Unit]
Description=Battery Input Manager

[Service]
Type=exec
ExecStart=@BIN_DIR@/bim

[Install]
WantedBy=default.target